# systemd/lightd.service
[Unit]
Description=LED strip Control Service

[Service]
Type=simple
ExecStart=/usr/bin/lightd

[Install]
WantedBy=multi-user.target

